# ansible/bsmiweb.service.j2
[Unit]
Description=bsmiweb web server
After=network.target

[Service]
Type=simple
WorkingDirectory={{ deploy_path }}
Environment=NODE_ENV=production
EnvironmentFile={{ deploy_path }}/.env
ExecStart={{ mise_node_bin }}/node src/index.js
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
